# user_mode/roms/skiboot/external/opal-prd/opal-prd.socket
[Unit]
Description=OPAL PRD Listen Socket

[Socket]
ListenStream=/run/opal-prd-control

[Install]
WantedBy=sockets.target
